# deploy/central/systemd/beakmeshwall-central.service
[Unit]
Description=BeakMeshWall Central Server
Documentation=https://github.com/your-org/beakmeshwall
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=beakmeshwall
Group=beakmeshwall
WorkingDirectory=/opt/beakmeshwall/central

# Load secrets and environment overrides from this file
EnvironmentFile=/etc/beakmeshwall/central.env

ExecStart=/opt/beakmeshwall/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    "app:create_app()"

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/beakmeshwall
PrivateTmp=true

[Install]
WantedBy=multi-user.target
